# target/linux/ipq40xx/files/arch/arm/boot/dts/files/etc/uci-defaults/seed.config
CONFIG_TARGET_ipq40xx=y
CONFIG_TARGET_ipq40xx_generic=y
CONFIG_TARGET_DEVICE_ipq40xx_generic_DEVICE_davolink_gapd-7200=y

# واجهة التحكم الرسومية المتقدمة LuCI
CONFIG_PACKAGE_luci=y
CONFIG_PACKAGE_luci-ssl=y

# دمج الميزات المطلوبة للشبكة
CONFIG_PACKAGE_kmod-eoip=y
CONFIG_PACKAGE_lldpd=y
CONFIG_PACKAGE_luci-app-watchcat=y
CONFIG_PACKAGE_luci-app-ttyd=y

